# open-ports.service
[Service]
Type=oneshot
ExecStart=/usr/bin/open_ports.sh
